# scienceclaw-heartbeat.service
[Unit]
Description=ScienceClaw Heartbeat Daemon — Delphi Mission
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/user/aether
ExecStart=/usr/bin/python3 /home/user/aether/heartbeat.py --interval 300
Restart=always
RestartSec=60
StandardOutput=append:/home/user/.scienceclaw/heartbeat_daemon.log
StandardError=append:/home/user/.scienceclaw/heartbeat_daemon.log
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/user"
Environment="PYTHONUNBUFFERED=1"
# Delphi mission environment
Environment="SCIENCECLAW_MISSION=delphi-qpi-2026"
Environment="SCIENCECLAW_DISABLE_SKILL_AUTO_REFRESH=1"

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/user/.scienceclaw /home/user/aether
PrivateTmp=true

[Install]
WantedBy=default.target
